# iml-update-handler.socket
[Socket]
ListenStream=/var/run/iml-update-handler.sock
RemoveOnStop=true

[Install]
WantedBy=sockets.target
